# nailgun-agent.cron
* * * * * root flock -w 0 -o /var/lock/nailgun-agent.lock -c "/usr/bin/nailgun-agent > /dev/null 2>&1"
